# data/supervisor.conf
[supervisord]
logfile=/home/octoprint/.octoprint/supervisord.log
stdout_logfile=/home/octoprint/.octoprint/%(program_name)s-stdout.log
stderr_logfile=/home/octoprint/.octoprint/%(program_name)s-stderr.log
nodaemon=true

[program:socat]
user=octoprint
command=/bin/bash -c "socat pty,wait-slave,link=/tmp/printer,perm=0777 tcp:klipper:4001,forever,interval=10"
redirect_stderr=true
redirect_stderr

[program:octoprint]
user=octoprint
command=/bin/bash -c "octoprint serve"
redirect_stderr=true